# tools/docker/arch/mnt/etc/systemd/system/docker-registry.service
[Unit]
Description=DOCKER-REGISTRY

[Service]
Type=simple
User=deploy
ExecStart=/var/www/itpkg/shared/scripts/docker-registry.sh

[Install]
WantedBy=multi-user.target

